# poo.conf
description "Poo/Web MAApi"
author "go akahane <akahane@golaboratory.jp>"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn

exec /poo/poo.bin

